# conf/supervisor.conf
[program:demo]
command = /www/demo/start_server
user = root
